# systemd/logos-llm.service
[Unit]
Description=Local LLM (llama.cpp server, Qwen3-1.7B) for the field node's "Local AI" chat
After=network.target

[Service]
ExecStart=/home/vpavlin/llama.cpp/build/bin/llama-server -m /home/vpavlin/models/Qwen3-1.7B-Q4_K_M.gguf --host 0.0.0.0 --port 8081 -c 2048 -t 3 --reasoning off
Restart=on-failure

[Install]
WantedBy=default.target
