# systemd/mesh-kiosk.service
[Unit]
Description=Basecamp Mesh Gateway kiosk (cage + chromium on SPI panel)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mesh-kiosk.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
